# game-server-restart.timer
[Timer]
Unit=%i-restart.service
OnCalendar=00:00 US/Central

[Install]
WantedBy=timers.target
